# ansible/mcp-session-purge.timer
[Unit]
Description=Daily purge of expired SuiteCRM MCP sessions

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
